# root/etc/services.d/uwsgi/run
#!/usr/bin/with-contenv sh

exec uwsgi --ini /etc/uwsgi.d/avss.ini;
