# root/etc/services.d/certbot/run
#!/usr/bin/with-contenv sh

exec certbot --nginx -d ${DOMAIN} --test-cert --agree-tos -n -m ${MAIL}
