# root/etc/periodic/daily/certbot
#!/usr/bin/with-contenv sh

exec certbot --nginx -d ${DOMAIN} --agree-tos -n -m ${MAIL};
